# deploy/systemd/mysql.service
[Unit]
Description=The MySQL Database Server
After=network.target

[Service]
Type=simple
User=mysql
Group=mysql
ExecStart=/usr/sbin/mysqld
ExecReload=/usr/sbin/mysqld shutdown
ExecStop=/usr/sbin/mysqld shutdown

[Install]
WantedBy=multi-user.target